# scripts/.local/bin/bk-manager/install/bk-manager.service
[Unit]
Description=BK Manager Scheduler

[Service]
Type=oneshot
ExecStart=/usr/local/bin/bk-manager run --scheduled
